# kmonad/.config/systemd/user/kmonad.service
[Unit]
Description=Kmonad Keyboard Remapping
After=graphical.target

[Service]
ExecStart=/home/drew/.local/bin/kmonad /home/drew/.config/kmonad/config.kbd
Restart=on-failure

[Install]
WantedBy=default.target
